# mkosi/mkosi.extra/etc/systemd/system/router-backup.timer
[Unit]
Description=Run a backup of user data once a day

[Timer]
OnActiveSec=1d
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target